# websocket_server.service
[Unit]
Description=WebSocket Server for M17 Reflector and DMR Dashboards
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/digidash/websocket_server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

